# drjimbo-fcgi.service
[Unit]
Description=drjimbo-game

[Service]
ExecStart=/bin/bash /opt/drjimbo-game/server.sh
User=www-data

[Install]
WantedBy=multi-user.target
